# roles/workstation/files/etc/profile.d/mac.sh
#!/bin/bash

source $(brew --prefix)/etc/bash_completion
